# templates/stacks/watchtower/docker.watchtower.service.j2
[Unit]
Description=Docker compose watchtower app
PartOf=docker.service
After=docker.service

[Service]
Type=oneshot
RemainAfterExit=true
WorkingDirectory=/opt/core/apps/watchtower
ExecStart=/usr/bin/docker compose up -d --remove-orphans
ExecStop=/usr/bin/docker compose down
StandardOutput=syslog
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target